# .config/systemd/user/m2i.service
[Unit]
Description=m2i service

[Service]
ExecStart=/usr/bin/m2i

[Install]
WantedBy=xdg-desktop-autostart.target
